# files/systemd/gunicorn.socket
[Unit]
Description={{project_name}} gunicorn socket

[Socket]
ListenStream=/run/{{project_name}}_gunicorn/socket

[Install]
WantedBy=sockets.target
